# claude-code-monitor.service
[Unit]
Description=Claude Code Monitor
After=network.target

[Service]
Type=simple
User=paulo
WorkingDirectory=/home/paulo/Github/claude-code-monitor
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
